# pdf-font-decoder.service
[Service]
User=web
WorkingDirectory=/srv/apps/pdf-font-decoder
ExecStart=mise run api
Restart=on-failure
RestartSec=5s
StandardError=append:/var/log/pdf-font-decoder-error.log

[Install]
WantedBy=multi-user.target
